# src/main/conf/ccm-purge.cron
# run every day at 00:30 ccm-purge
#

30 0 * * * root /usr/sbin/ccm-purge >>/var/log/ccm-purge.log 2>&1

