# src/main/conf/ccm-fetch.cron
# To avoid lost notifications, run ccm-fetch regularly.
#
RANDOM_DELAY=60
0 * * * * root (date --iso-8601=seconds --utc; /usr/sbin/ccm-fetch) >>/var/log/ccm-fetch.log 2>&1
